# deploy/systemd/flipdisc-worker@.service
[Unit]
Description=Flip-Disc Media Pipeline %i (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# These values are illustrative; use deploy/install_systemd.sh to install with correct paths.
User=pi
Group=pi
WorkingDirectory=/home/pi/flip-disc/media_pipeline
EnvironmentFile=/etc/flipdisc/worker-%i.env
Environment=PYTHONPATH=/home/pi/flip-disc
ExecStart=/home/pi/flip-disc/media_pipeline/.venv/bin/uv run python media_pipeline/runner.py %i
Restart=on-failure
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
