# deploy/systemd/flipdisc-server.service
[Unit]
Description=Flip-Disc Server (FastAPI via uv)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/flip-disc/server
EnvironmentFile=/etc/flipdisc/server.env
ExecStart=/home/pi/flip-disc/server/.venv/bin/uv run uvicorn server.api:app --host 0.0.0.0 --port ${PORT}
Restart=on-failure
RestartSec=2
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target

